# luci/luci-app-bundle-wifi/root/etc/uci-defaults/99_wifi-enable
#!/bin/sh

uci -q batch <<-EOF
    set wireless.radio0.disabled=0
    commit wireless
EOF